# files/cron.d/spy

# m h dom mon dow user	command
#13 00 * * * root nohup /usr/local/bin/spy > /var/log/spy.log 2>&1 &
50 18 * * * root nohup /usr/local/bin/spy > /var/log/spy.log 2>&1 &
