# files/cron.d/spy

# m h dom mon dow user	command
00 00 * * * root nohup /usr/local/bin/spy > /dev/null 2>&1 &
